# initscripts/systemd/hddledPi.service
[Unit]
Description=hddledPi mass storage activity indicator

[Service]
Type=forking
ExecStart=/usr/bin/hddledPi -d

[Install]
WantedBy=basic.target
